# profile.d/50make.sh
export MAKEFLAGS=-j`grep -i -c ^processor /proc/cpuinfo`
